# templates/etc/profile.d/node-n.sh.j2
export N_PREFIX={{node_n_prefix}}
export PATH=$N_PREFIX/bin:$PATH
