# files/etc/profile.d/node-n.sh
export PATH=/opt/node-n/bin:$PATH
export N_PREFIX=/opt/node-n
