# files/etc/profile.d/node-n.sh
export N_PREFIX=/opt/node-n
export PATH=$N_PREFIX/bin:$PATH
